# build/bahmni-func-tests/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/etc/init.d/mysqld start
priority=3
user=root
