# build/openerp/openerp.conf
[program:openerp]
command = /usr/bin/openerp-server --update all --without-demo all --addons /opt/openerp/addons --db_user openerp --db_host %(ENV_DB_HOST)s
user = openerp
environment = USER=openerp,LOGNAME=openerp,HOME=/home/openerp
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/openerp-server.log